# .config/autostart/gmixer-trayicon.desktop
[Desktop Entry]
Hidden=true
Name=GMixer
Exec=gmixer -d
Type=Application

