# .config/autostart/ssh-add.desktop
[Desktop Entry]
Type=Application
Name=SSH Add
Exec=/usr/bin/ssh-add

